# iris-service.service
[Unit]
Description=Iris Recognition Service
Documentation=file:///moutum/iris-service/README.md
After=network.target

[Service]
Type=simple
User=root
Group=root

# 工作目录必须是native-libs/linux-x64
WorkingDirectory=/moutum/iris-service/native-libs/linux-x64

# 设置环境变量
Environment="LD_LIBRARY_PATH=/moutum/iris-service/native-libs/linux-x64"
Environment="LD_PRELOAD=/moutum/iris-service/native-libs/linux-x64/libirisengine.so:/moutum/iris-service/native-libs/linux-x64/libsbklic.so:/moutum/iris-service/native-libs/linux-x64/libsbkutil.so"
Environment="JAVA_HOME=/usr/lib/jvm/java-8-openjdk-amd64"

# 启动命令
ExecStart=/usr/bin/java -Djava.library.path=. \
    -Xmx512m \
    -Xms256m \
    -jar /moutum/iris-service/iris-springboot-1.0.0.jar

# 重启策略
Restart=always
RestartSec=10

# 标准输出和错误输出
StandardOutput=journal
StandardError=journal

# 资源限制
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target

# ========================================
# 安装说明：
# ========================================
# 
# 1. 根据实际部署路径修改：
#    - WorkingDirectory: /moutum/iris-service/native-libs/linux-x64
#    - LD_LIBRARY_PATH: /moutum/iris-service/native-libs/linux-x64
#    - JAVA_HOME: 实际JDK路径（如需要）
#    - ExecStart中的jar路径: /moutum/iris-service/iris-springboot-1.0.0.jar
#
# 2. 复制到systemd目录：
#    sudo cp iris-service.service /etc/systemd/system/
#
# 3. 重新加载systemd：
#    sudo systemctl daemon-reload
#
# 4. 启用开机自启：
#    sudo systemctl enable iris-service
#
# 5. 启动服务：
#    sudo systemctl start iris-service
#
# 6. 查看状态：
#    sudo systemctl status iris-service
#
# 7. 查看日志：
#    sudo journalctl -u iris-service -f
#
# 8. 停止服务：
#    sudo systemctl stop iris-service
#
# 9. 重启服务：
#    sudo systemctl restart iris-service
#

